# .config/systemd/user/wallpaper.timer
[Unit]
Description=Run hourly job

[Timer]
OnCalendar=*:00

[Install]
WantedBy=default.target
; WantedBy=graphical-session.target
